# shw/roles/jenkins/templates/jenkins.service
[Unit]
Description=Jenkins Daemon

[Service]
ExecStart=/data/services/startjenkins.sh
User=root

[Install]
WantedBy=multi-user.target
