# shw/roles/rocketmq/templates/rocketmq-console.ini
[program:rocketmq-console]
command=java -jar rocketmq-console-ng-1.0.0.jar
numprocs=1
directory=/data/softwares/rocketmq-console
autostart=true
autorestart=true
startretries=3
user=root
stdout_logfile=/data/logs/rocketmq-console/stdout.log
stopsignal=KILL
redirect_stderr=true
stopasgroup=true
